# src/systemd/ossec-server.target
[Unit]
Description=The OSSEC HIDS server
After=network.target 

[Install]
Wanted=ossec-dbd.service